# roles/udev/files/automount-by-label.rules
KERNEL=="sd[a-z][0-9]", SUBSYSTEMS=="usb", ACTION=="add", RUN+="/usr/bin/usbmount.sh --mount %E{DEVNAME} %E{ID_FS_LABEL}"
